# initramfs/etc2/init.d/ammit-etc
#!/sbin/openrc-run
# Copyright 2020 Gentoo Authors
# Distributed under the terms of the GNU General Public License v2

name="etc subset updater"
description="saves changes in wanted items under /etc to persistent storage"
command=/sbin/ammit-etc.sh
command_background=true
pidfile=/var/run/ammit-etc.pid

depend() {
    after localmount
}
