# initramfs/etc2/local.d/nvidia.start
#!/bin/bash

depend() {
    need localmount
    before xdm
}

eselect opencl set nvidia
eselect opengl set nvidia
